# roles/cart-p/templates/systemd.service
[Unit]
Description = Cart Service
[Service]
User=roboshop
Environment=REDIS_HOST=redis-p.roboshop.internal
Environment=CATALOGUE_HOST=catalogue-p.roboshop.internal
ExecStart=/bin/node /home/roboshop/cart/server.js
SyslogIdentifier=cart

[Install]
WantedBy=multi-user.target